# profile.d/alpi.sh
#!/bin/sh
set -e

function alpi_install() {

}
