# init.d/alpi.sh
#!/sbin/openrc-run
command=echo "runned"
