# core/profile.d/alias.sh
#!/bin/sh
alias ..='cd ..'
alias ...='cd ../..'

#function nginxinit() {
  #adduser -D -g 'www' www
  #mkdir /www
  #chown -R www:www /var/lib/nginx
  #chown -R www:www /www
#}
#alias nginx-init='nginxinit'
#alias nginx-start='rc-service nginx start'
#alias nginx-enable='rc-update add nginx'

#alias squid-start='rc-service squid start'
#alias squid-enable='rc-update add squid'
#alias squid-reconfigure='squid -k reconfigure'
#alias squid-check='squid -k check'
#alias squid-stat='netstat -tl'


